# tidy/profile.d/tidy.sh
# Set a configuration file path for tidy(1)
HTML_TIDY=$HOME/.tidyrc
export HTML_TIDY
