# sh/profile.d/browser.sh
# Set command-line browser to lynx; ~/.xinitrc will change this to something
# graphical instead
BROWSER=lynx
export BROWSER
