# sh/profile.d/editor.sh
# Set command-line editor
EDITOR=ed
export EDITOR
