# sh/profile.d/games.sh
# Add ~/.local/games to PATH if it exists
[ -d "$HOME"/.local/games ] && PATH=$HOME/.local/games:$PATH
