# sh/profile.d/go.sh
# Define path for Go code to be installed
GOPATH=$HOME/.local/gocode
export GOPATH

# Prepend GOPATH to PATH for the executables
PATH=$GOPATH/bin:$PATH

