# sh/profile.d/timezone.sh
# Timezone
# TZ=Pacific/Auckland
TZ=America/Montreal
export TZ
