# sh/profile.d/visual.sh
# For a visual editor, use whichever kind of vi we get when we invoke 'vi'
VISUAL='vi'
export VISUAL
