# sh/profile.d/visual.sh
# Use first found implementation of vi(1)
VISUAL='vi'
export VISUAL
