# sh/profile.d/visual.sh
# Set visual editor
VISUAL=vi
export VISUAL
